# scripts/services/github-mirror.timer
[Unit]
Description=Periodic GitHub mirror sync

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
